# supervisod.conf
[supervisord]
logfile=supervisord.log                         ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/home/ec2-user/second/logs/                ; where child log files will live

[program:second]
directory=/home/ec2-user/second/
command=java -jar second-0.0.1-SNAPSHOT.jar --spring.profiles.active=production
stdout_logfile=syslog
stderr_logfile=syslog
user=ec2-user
